# deploy/ansible/roles/relay_agent/templates/relay-agent.service.j2
[Unit]
Description=Relay Agent — HTTP sidecar for xray-relay status
After=network.target xray-relay.service
Wants=xray-relay.service

[Service]
Type=simple
User={{ relay_agent_user }}
WorkingDirectory={{ relay_agent_install_dir }}
ExecStart={{ relay_agent_install_dir }}/venv/bin/python {{ relay_agent_install_dir }}/main.py
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
